# provisioning/roles/wiki/templates/supervisor.conf.j2
[program: wiki]
command={{ cmd }}
user={{ user}}
autostart=True
directory={{ workdir }}
stdout_logfile={{ log_file }}
stderr_logfile={{ log_file }}